# systemd/clickhouse-backup.service
[Unit]
Description=ClickHouse → S3 Backup (clickhouse-backup v1.x)
After=network-online.target clickhouse-server.service
Wants=network-online.target

[Service]
Type=oneshot
User=clickhouse
Group=clickhouse
EnvironmentFile=/etc/clickhouse-backup/backup.env

# clickhouse-backup читает данные напрямую с диска — запускаем от имени пользователя clickhouse
ExecStart=/opt/clickhouse-backup/venv/bin/python /opt/clickhouse-backup/backup_clickhouse.py

# Ограничение ресурсов
TimeoutStartSec=7200
MemoryMax=1G

# Логи: journalctl -u clickhouse-backup -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clickhouse-backup

# Не перезапускать автоматически при ошибке (ошибки бэкапа требуют ручного разбора)
Restart=no

[Install]
WantedBy=multi-user.target
